# templates/supervisord-jarn.conf
[unix_http_server]
file=${buildout:directory}/var/supervisor.sock
chmod=0600

[supervisorctl]
serverurl=unix://${buildout:directory}/var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=${buildout:directory}/var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=${buildout:directory}/var/supervisord.pid ;
childlogdir=${buildout:directory}/var/log
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
directory=%(here)s

[program:memcached]
command = %(here)s/../bin/memcached -m 64 -p 11211 -l 127.0.0.1
environment = LD_LIBRARY_PATH=%(here)s/../parts/libevent/lib/
autostart = true
autorestart = true
priority = 100

[program:haproxy]
command = ${buildout:bin-directory}/haproxy -f %(here)s/haproxy.cfg -db
autostart= true
autorestart = true
redirect_stderr = true
priority = 101

[program:varnish]
command = ${buildout:bin-directory}/varnish
startretries = 5
redirect_stderr = true
autostart = true
autorestart = true
priority = 102

[program:ejabberd]
command = ${buildout:bin-directory}/ejabberd
startretries = 5
redirect_stderr = true
autostart = true
autorestart = true
priority = 103

[program:instance1]
command = ${buildout:bin-directory}/instance1 console
startretries = 5
redirect_stderr = true
autostart= true
autorestart = true
priority = 201

[program:instance2]
command = ${buildout:bin-directory}/instance2 console
startretries = 5
redirect_stderr = true
autostart= true
autorestart = true
priority = 202

[group:zope]
programs = instance1,instance2
priority = 210

[eventlistener:memmon]
command=${buildout:bin-directory}/memmon -g zope=250MB
events=TICK_60
